# homeassistant/loxone2velux/rootfs/etc/services.d/loxone2velux/finish
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Loxone2Velux Gateway - s6 service finish script
# ==============================================================================

declare exit_code=${1}

if [[ "${exit_code}" -ne 0 ]] && [[ "${exit_code}" -ne 256 ]]; then
    bashio::log.warning "Loxone2Velux Gateway crashed with exit code ${exit_code}"
    echo "${exit_code}" > /run/s6-linux-init-container-results/exitcode
    exec /run/s6/basedir/bin/halt
fi

bashio::log.info "Loxone2Velux Gateway stopped, restarting..."
